# rootcopy/usr/share/systemd/system/bluefora-firstrun.service
[Service]
ExecStart=bash /usr/lib/bluefora/firstrun.sh
ConditionFirstBoot=yes